# supervisord.conf
[unix_http_server]
file=supervisord.sock

[inet_http_server]
port=:9001

[supervisord]
logfile=supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://supervisord.sock

[program:cendyneyells]
command=python3 yellbot.py
redirect_stderr=true
stdout_logfile=yell.log
autostart=true

[program:cendynesays]
command=python3 bot.py
redirect_stderr=true
stdout_logfile=says.log
autostart=true
